# mfc/nokeena/src/install_files/nkn_sys_mon.cron
# run the nkn_sys_mon.sh every one minute
*/1 * * * * root /opt/nkn/bin/nkn_sys_mon.sh
MAILTO=""
