# samara/tree/src/base_os/linux_el/el6/initscripts_files/serial.conf
# ttyS0 - agetty
#
# This service maintains an agetty on ttyS0

start on stopped rc RUNLEVEL=[345]
stop on starting rc RUNLEVEL=[!345]

respawn
exec /sbin/agetty -h ttyS0 9600 vt102
